# raspberrypi/crontab
*/5	*	* * *	~/LGSStatus/raspberrypi/stats.sh
